# tests/containers/entrypoint/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
childlogdir = /var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
loglevel=debug

[include]
files=/etc/supervisord.d/*.conf

[program:condor-ce]
command=/usr/share/condor-ce/condor_ce_startup -f
autorestart=true
startsecs=20

[program:condor]
user = root
command = /usr/sbin/condor_master -f
autorestart=true
startsecs=20

[program:munge]
user = munge
command = /usr/sbin/munged -F
autorestart = true
startsecs = 20

[program:mariadb]
user = mysql
command = /usr/bin/mysqld_safe --basedir=/usr --init-file=/etc/condor-ce/image-init.d/init.sql
priority = 50
autorestart = true
startsecs = 20

[program:slurmdbd]
user = root
command = /usr/sbin/slurmdbd -D
priority = 60
autorestart = true
startsecs = 20

[program:slurmd]
user = root
command = /usr/sbin/slurmd -D
autorestart = true
startsecs = 20

[program:slurmctld]
user = root
command = /usr/sbin/slurmctld -D
autorestart = true
startsecs = 20
